# emf-corosync.service
[Unit]
Description=EMF Corosync Service
PartOf=emf-manager.target
After=rabbitmq-server.service

[Service]
Type=simple
EnvironmentFile=-/var/lib/chroma/overrides.conf
ExecStart=/bin/python /usr/share/chroma-manager/manage.py chroma_service --name=corosync corosync
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal