# device-scanner/mount-emitter/systemd-units/swap-emitter.service
[Unit]
Description=IML Swap Emitter

[Service]
ExecStart=/bin/bash -c 'exec /usr/bin/findmnt -P -s -e -t swap | /usr/bin/mount-emitter'
Type=oneshot
