# kuma/systemd-units/emf-postgres-sidecar.service
[Unit]
Description=EMF Postgres sidecar
Before=postgresql-13.service
Requires=postgresql-13.service
PartOf=emf-manager.target

[Service]
EnvironmentFile=/etc/emf/bootstrap.conf
EnvironmentFile=-/etc/emf/overrides.conf
ExecStartPre=/bin/rm -f /tmp/kuma-al-postgres-default.sock
ExecStart=/usr/bin/kuma-dp run \
    --cp-address=https://${CP_ADDR}:5678 \
    --dataplane-file=/etc/emf/dataplanes/postgres.yml \
    --dataplane-token-file=/etc/emf/tokens/dataplane-token \
    --dataplane-var postgres_port=${PG_PORT} \
    --dataplane-var mgmt_address=${MGMT_ADDR} \
    --admin-port=
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
LimitNOFILE=4096
