# kuma/systemd-units/emf-grafana-sidecar.service
[Unit]
Description=EMF Grafana sidecar
Requires=grafana-server.service
PartOf=grafana-server.service
PartOf=emf-manager.target

[Service]
EnvironmentFile=/etc/emf/bootstrap.conf
EnvironmentFile=-/etc/emf/overrides.conf
ExecStartPre=/bin/rm -f /tmp/kuma-al-grafana-default.sock
ExecStart=/usr/bin/kuma-dp run \
    --cp-address=https://${CP_ADDR}:5678 \
    --dataplane-file=/etc/emf/dataplanes/grafana.yml \
    --dataplane-token-file=/etc/emf/tokens/dataplane-token \
    --dataplane-var grafana_port=${GRAFANA_PORT} \
    --dataplane-var mgmt_address=${MGMT_ADDR} \
    --admin-port=
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
LimitNOFILE=4096
