# chroma-agent/chroma-agent.service
[Unit]
Description=IML Agent Daemon
After=network.target

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/chroma-agent-daemon
Restart=on-failure
StandardOutput=journal
StandardError=journal