# iml-agent/systemd-units/rust-iml-agent.path
[Unit]
Description=Rust IML Agent path detection
PartOf=iml-storage-server.target

[Path]
PathExists=/etc/iml/manager-url.conf

[Install]
WantedBy=iml-storage-server.target
