# emf-services/systemd-units/emf-host.service
[Unit]
Description=EMF Host Service
PartOf=emf-manager.target
Wants=emf-host-sidecar.service
After=emf-host-sidecar.service

[Service]
Type=simple
Environment=RUST_LOG=info,sqlx::query=warn
EnvironmentFile=/etc/emf/bootstrap.conf
EnvironmentFile=-/etc/emf/overrides.conf
ExecStart=/bin/emf-host
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=emf-manager.target
