# emf-services/systemd-units/emf-rust-corosync.service
[Unit]
Description=EMF Rust Corosync Service
PartOf=emf-manager.target
Wants=emf-corosync-sidecar.service
After=emf-corosync-sidecar.service

[Service]
Type=simple
Environment=RUST_LOG=info,sqlx::query=warn
EnvironmentFile=/etc/emf/bootstrap.conf
EnvironmentFile=-/etc/emf/overrides.conf
ExecStart=/bin/emf-corosync
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=emf-manager.target
