# emf-services/systemd-units/emf-device.service
[Unit]
Description=EMF Device Service
PartOf=emf-manager.target
Wants=emf-device-sidecar.service
After=emf-device-sidecar.service

[Service]
Type=simple
Environment=RUST_LOG=info,sqlx::query=warn
EnvironmentFile=/etc/emf/bootstrap.conf
EnvironmentFile=-/etc/emf/overrides.conf
ExecStart=/bin/emf-device
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=emf-manager.target
